# templates/etc/profile.d/vault.sh.j2
#!/usr/bin/env bash
export VAULT_ADDR={{ vault_addr }}
export VAULT_CACERT={{ vault_cacert }}